# www/bitrix/crontab/crontab.cfg




* * * * * php -f /home/b2153/domains/itprom.ru/www/bitrix/modules/main/tools/cron_events.php
